# debian/mcp-bridge.service
[Unit]
Description=MCP Bridge
Documentation=https://github.com/TheParadox20/server-mcp-bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mcp
Group=mcp
WorkingDirectory=/var/lib/mcp_bridge
ExecStart=/usr/bin/mcp_bridge daemon --config /etc/mcp_bridge/mcp.json
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5

# Security hardening
#
# NoNewPrivileges is intentionally NOT set: the daemon spawns the setuid
# helper /usr/lib/mcp_bridge/mcp_bridge-priv (mode 4750 root:mcp) for
# privileged operations (useradd, install/revoke sudoers drop-ins), and
# tool calls drop privileges to per-session mcp_user_<shortid> accounts
# which then invoke sudo against grants in /etc/sudoers.d. NoNewPrivileges
# would mask the setuid bit on the helper *and* on sudo, breaking both.
#
# AmbientCapabilities gives the daemon (running as `mcp`) the ability to
# setuid/setgid down to the per-session user before exec'ing tool commands.
# CapabilityBoundingSet caps what the daemon can ever hold, even via the
# setuid helper, to the same minimum.
AmbientCapabilities=CAP_SETUID CAP_SETGID
CapabilityBoundingSet=CAP_SETUID CAP_SETGID CAP_CHOWN CAP_FOWNER CAP_DAC_OVERRIDE CAP_DAC_READ_SEARCH

ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/mcp_bridge /var/log/mcp_bridge /etc/sudoers.d
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mcp_bridge

[Install]
WantedBy=multi-user.target
