# debian/postinst
#!/bin/sh
set -e

CONF=/etc/mcp_bridge/mcp.json
TEMPLATE=/usr/share/mcp_bridge/mcp.json.template

case "$1" in
    configure)
        # Create service user if it doesn't exist
        if ! getent passwd mcp >/dev/null 2>&1; then
            adduser --system --group --no-create-home \
                --home /var/lib/mcp_bridge \
                mcp
        fi

        # Ensure directories exist with the right ownership and modes
        install -d -m 750 -o root -g mcp /etc/mcp_bridge
        install -d -m 750 -o mcp  -g mcp /var/lib/mcp_bridge
        install -d -m 750 -o mcp  -g mcp /var/lib/mcp_bridge/users
        install -d -m 750 -o mcp  -g mcp /var/lib/mcp_bridge/state
        install -d -m 750 -o mcp  -g mcp /var/log/mcp_bridge

        # Materialize the JSON config from the template if not already present
        if [ ! -f "$CONF" ]; then
            if ! command -v openssl >/dev/null 2>&1; then
                echo "ERROR: openssl is required to generate the admin token." >&2
                exit 1
            fi
            if ! command -v sha256sum >/dev/null 2>&1; then
                echo "ERROR: sha256sum is required to generate the admin token hash." >&2
                exit 1
            fi

            SALT=$(openssl rand -hex 32)
            TOKEN=$(openssl rand -hex 32)
            HASH=$(printf '%s%s' "$SALT" "$TOKEN" | sha256sum | awk '{print $1}')

            umask 0177
            sed -e "s|__GLOBAL_TOKEN_SALT__|$SALT|" \
                -e "s|__ADMIN_TOKEN_HASH__|$HASH|" \
                "$TEMPLATE" > "$CONF.tmp"
            mv "$CONF.tmp" "$CONF"
            chown root:mcp "$CONF"
            chmod 0640 "$CONF"

            echo ""
            echo "=========================================="
            echo " MCP Bridge installed!"
            echo "=========================================="
            echo ""
            echo " Config:    $CONF"
            echo " Endpoint:  http://localhost:8080/"
            echo ""
            echo " Admin API token (save this — it is shown only once):"
            echo "   $TOKEN"
            echo ""
            echo " To start the service:"
            echo "   sudo systemctl enable --now mcp-bridge"
            echo ""
            echo "=========================================="
        fi

        # Setuid helper: ownership root:mcp, mode 4750
        if [ -x /usr/lib/mcp_bridge/mcp_bridge-priv ]; then
            chown root:mcp /usr/lib/mcp_bridge/mcp_bridge-priv
            chmod 4750     /usr/lib/mcp_bridge/mcp_bridge-priv

            # Provision the ephemeral OS account that the install admin token
            # binds to (see auth.cpp), and install its wildcard sudoers entry.
            # Both subcommands are idempotent across re-runs of postinst.
            #   useradd rc 14 = "user already exists" -> treat as success.
            set +e
            /usr/lib/mcp_bridge/mcp_bridge-priv useradd mcpadmin
            rc=$?
            set -e
            if [ "$rc" != 0 ] && [ "$rc" != 14 ]; then
                echo "WARNING: provisioning mcp_user_mcpadmin failed (rc=$rc);" >&2
                echo "         the install admin token will have no sudo authority." >&2
                echo "         Re-run: sudo /usr/lib/mcp_bridge/mcp_bridge-priv useradd mcpadmin" >&2
            else
                if ! /usr/lib/mcp_bridge/mcp_bridge-priv install-system-admin mcpadmin; then
                    echo "WARNING: install-system-admin failed; admin sudo not configured." >&2
                fi
            fi
        fi

        # Non-secret motd (admin token NEVER goes here — it's printed once on /dev/tty above).
        install -d -m 755 /etc/motd.d
        cat > /etc/motd.d/mcp_bridge <<'MOTD'

      ******         ******
    **      **     **      **
   *          *****          *
   *           ***           *
   *          *****          *
    **      **     **      **
      ******         ******

  MCP Bridge — running as system user `mcp` (not in the sudo group).

  config:    /etc/mcp_bridge/mcp.json
  state:     /var/lib/mcp_bridge/
  logs:      journalctl -u mcp-bridge   |   /var/log/mcp_bridge/server.log
  endpoint:  http://localhost:8080/
  service:   sudo systemctl status mcp-bridge

  Provision a new user:   sudo mcp_bridge auth create --name X --email Y
  Rotate a user's token:  sudo mcp_bridge auth rotate <shortid>

MOTD

        # Reload systemd
        if [ -d /run/systemd/system ]; then
            systemctl daemon-reload || true
        fi
        ;;
esac

#DEBHELPER#

exit 0
