# debian/postrm
#!/bin/sh
set -e

case "$1" in
    purge)
        # Remove sudoers drop-ins issued by the daemon and at install time
        rm -f /etc/sudoers.d/mcp_grant_*
        rm -f /etc/sudoers.d/mcp_system_admin

        # Remove the install-time admin's ephemeral OS account
        if getent passwd mcp_user_mcpadmin >/dev/null 2>&1; then
            deluser --remove-home mcp_user_mcpadmin || true
        fi

        # Remove config, logs, state, motd
        rm -rf /etc/mcp_bridge
        rm -rf /var/lib/mcp_bridge
        rm -rf /var/log/mcp_bridge
        rm -f  /etc/motd.d/mcp_bridge

        # Remove service user
        if getent passwd mcp >/dev/null 2>&1; then
            deluser --system mcp || true
        fi
        if getent group mcp >/dev/null 2>&1; then
            delgroup --system mcp || true
        fi
        ;;

    remove)
        # Stop service if running
        if [ -d /run/systemd/system ]; then
            systemctl stop mcp-bridge || true
            systemctl disable mcp-bridge || true
        fi
        ;;
esac

#DEBHELPER#

exit 0
